# ansible/roles/dbadmin/templates/dbadmin-proxy.service.j2
[Unit]
Description=DB Admin Auth Proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/default/dbadmin-proxy
ExecStart=/opt/hosting/bin/dbadmin-proxy
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
